# build/scripts/monotest.conf
description "MonoTest job"

start on startup

script
  mono /home/ubuntu/MonoTest/MonoTest.exe
end script
